# systemd/cloudflare-ddns.timer
[Unit]
Description=Run the Cloudflare DDNS script hourly

[Timer]
OnBootSec=1min
OnUnitActiveSec=1hr

[Install]
WantedBy=timers.target
